# scripts/dramatic-options.service
[Unit]
Description=dramatic_options Web App
After=network.target

[Service]
Type=simple
# DEV runs as the login user; PROD typically uses a dedicated service account
# and a /opt path (mirror real_options' coffee-bot / /opt/real_options if desired).
User=rodrigo
Group=rodrigo
WorkingDirectory=/home/rodrigo/dramatic_options
EnvironmentFile=/home/rodrigo/dramatic_options/.env
Environment="PYTHONUNBUFFERED=1"

# TODO: set the real entry point once the app exists. Examples:
#   uvicorn:   ExecStart=/home/rodrigo/dramatic_options/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
#   gunicorn:  ExecStart=/home/rodrigo/dramatic_options/venv/bin/gunicorn app:app -b 0.0.0.0:8000
#   streamlit: ExecStart=/home/rodrigo/dramatic_options/venv/bin/streamlit run app.py --server.address 0.0.0.0 --server.port 8501
ExecStart=/home/rodrigo/dramatic_options/venv/bin/python -u app.py

Restart=always
RestartSec=10

StandardOutput=append:/home/rodrigo/dramatic_options/logs/app.log
StandardError=append:/home/rodrigo/dramatic_options/logs/app.log

[Install]
WantedBy=multi-user.target
