# dashboard_web/deploy/dramatic-options-web.service
# Read-only observability web dashboard (React SPA + FastAPI /api, one port, tailnet-bound, keyless).
# Install: cp this to /etc/systemd/system/ then `systemctl daemon-reload && systemctl enable --now
# dramatic-options-web`. Parallel to dramatic-options-dashboard (Streamlit, :8502) — never the trade/auth path.
[Unit]
Description=Dramatic Options — observability web dashboard (read-only, keyless)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=rodrigo
WorkingDirectory=/home/rodrigo/dramatic_options-claude
ExecStart=/home/rodrigo/dramatic_options-claude/dashboard_web/deploy/run.sh
# on-failure (NOT always): a clean stop stays stopped; transient boot races retry slowly (the script polls
# the tailnet IP ~60s and exits fail-closed if absent, so this can't tight-loop against a missing IP).
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
